# deploy/systemd/fault-detection.service.example
[Unit]
Description=Industrial Fault Detection FastAPI Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/app/backend
EnvironmentFile=/opt/app/.env
ExecStart=/opt/app/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
